# config/ai.openclaw.self-heal.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.openclaw.self-heal</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/mac/.openclaw/workspace/scripts/openclaw-self-heal.sh</string>
  </array>

  <key>StartInterval</key>
  <integer>1800</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>WorkingDirectory</key>
  <string>/Users/mac/.openclaw/workspace</string>

  <key>StandardOutPath</key>
  <string>/Users/mac/.openclaw/workspace/state/guard/launchd-self-heal.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/mac/.openclaw/workspace/state/guard/launchd-self-heal.err.log</string>
</dict>
</plist>
